# files/system/usr/lib/systemd/system/flatpak-system-update.service
[Unit]
Description=Update System Flatpaks
Documentation=man:flatpak
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStartPre=/usr/bin/nm-online -q --timeout=60
ExecStart=/usr/bin/flatpak update -y --noninteractive --system

Nice=19
IOSchedulingClass=idle