# files/system/usr/lib/systemd/system/setup-kvm-once@.service
[Unit]
Description=Setup KVM/libvirt permissions for user %i
After=libvirtd.service
Wants=libvirtd.service

[Service]
Type=oneshot
ExecStart=/usr/bin/setup-kvm.sh %i
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
