# files/system/usr/lib/systemd/user/flatpak-user-update.service
[Unit]
Description=Update User Flatpaks
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStartPre=/usr/bin/nm-online -q --timeout=60
ExecStart=/usr/bin/flatpak update -y --noninteractive --user
Nice=19
IOSchedulingClass=idle