# files/system/usr/lib/systemd/user/podman-prune-user.timer
[Unit]
Description=Run rootless Podman prune daily

[Timer]
# Daily cleanup after boot stabilization; persistent to catch missed runs.
OnBootSec=25m
OnUnitActiveSec=1d
RandomizedDelaySec=30m
Persistent=true
AccuracySec=1m
Unit=podman-prune-user.service

[Install]
WantedBy=timers.target
