# files/system/usr/lib/systemd/user/rclone@.service
[Unit]
Description=Rclone Mount - %i
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
# Carrega as variáveis do provedor
EnvironmentFile=-%h/.config/rclone/env/%i.env

# Cria o ponto de montagem automaticamente
ExecStartPre=/usr/bin/mkdir -p %h/mnt/%i

# Execução nativa (sem --daemon) suportando notificações do systemd
ExecStart=/usr/bin/rclone mount %i: %h/mnt/%i \
  --vfs-cache-mode full \
  --vfs-cache-max-size 50G \
  --vfs-cache-max-age 24h \
  --dir-cache-time 1h \
  --log-file=%h/.local/log/rclone-%i.log \
  --log-level INFO \
  $RCLONE_FLAGS

# Fedora Kinoite utiliza fusermount3
ExecStop=/usr/bin/fusermount3 -u %h/mnt/%i

Restart=on-failure
RestartSec=10

# Limites adequados para 64GB de RAM e CPU de 32 threads
MemoryMax=4G
LimitNOFILE=1048576

[Install]
WantedBy=default.target
