# files/system/usr/lib/systemd/user/topgrade-system-update.timer
[Unit]
Description=Run topgrade system update

[Timer]
OnStartupSec=15m
OnUnitInactiveSec=2h
RandomizedDelaySec=5m
AccuracySec=1m
Persistent=true
Unit=topgrade-system-update.service

[Install]
WantedBy=timers.target
