# files/system/usr/lib/systemd/user/topgrade-boot-update.timer
[Unit]
Description=Run topgrade update (misc/containers)

[Timer]
OnStartupSec=30m
OnUnitInactiveSec=6h
RandomizedDelaySec=10m
AccuracySec=1m
Persistent=true
Unit=topgrade-boot-update.service

[Install]
WantedBy=timers.target
