# files/system/usr/lib/systemd/user/topgrade-update.timer
[Unit]
Description=Run topgrade update

[Timer]
OnStartupSec=5m
OnUnitInactiveSec=1h
RandomizedDelaySec=5m
AccuracySec=1m
Persistent=true
Unit=topgrade-update.service

[Install]
WantedBy=timers.target
