# files/system/usr/lib/systemd/user/topgrade-flatpak-update.timer
[Unit]
Description=Run topgrade flatpak update

[Timer]
# Start 5 minutes after user session begins, then every 6 hours
# Reduced frequency from 1h to 6h to avoid excessive network traffic
OnStartupSec=5m
OnUnitInactiveSec=6h
RandomizedDelaySec=30m
AccuracySec=1m
Persistent=true
Unit=topgrade-flatpak-update.service

[Install]
WantedBy=timers.target
