# files/system/usr/lib/systemd/user/topgrade-boot-update.timer
[Unit]
Description=Run topgrade system update shortly after boot

[Timer]
OnStartupSec=30m
OnUnitInactiveSec=6h
RandomizedDelaySec=15m
Persistent=true
Unit=topgrade-boot-update.service

[Install]
WantedBy=timers.target
