# files/system/usr/lib/systemd/user/topgrade-update.service
[Unit]
Description=Topgrade unattended update (system + flatpak)
Documentation=https://github.com/topgrade-rs/topgrade
Wants=network-online.target
After=network-online.target
StartLimitIntervalSec=15min
StartLimitBurst=2

[Service]
Type=oneshot
# flock exit code 75 indicates lock contention and is treated as success.
SuccessExitStatus=75
TimeoutStartSec=2h
ExecStart=/usr/bin/flock -E 75 -w 900 %t/topgrade-update.lock /usr/bin/topgrade -cy --no-ask-retry --auto-retry 2 --notify-end on_failure --only system flatpak
Nice=19
IOSchedulingClass=idle
CPUSchedulingPolicy=idle
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
