# files/system/usr/lib/systemd/user/topgrade-system-update.timer
[Unit]
Description=Run topgrade system update daily

[Timer]
OnStartupSec=30m
OnUnitInactiveSec=4h
RandomizedDelaySec=30m
Persistent=true
Unit=topgrade-system-update.service

[Install]
WantedBy=timers.target
