# files/system/usr/lib/systemd/user/topgrade-flatpak-update.timer
[Unit]
Description=Run topgrade flatpak update every 6 hours

[Timer]
OnStartupSec=10m
OnUnitInactiveSec=2h
RandomizedDelaySec=10m
Persistent=true
Unit=topgrade-flatpak-update.service

[Install]
WantedBy=timers.target
