# files/system/usr/lib/systemd/user/topgrade-system-update.service
[Unit]
Description=Topgrade update (system only)
Wants=network-online.target
After=network-online.target
# ConditionACPower=true removed for desktop workstation compatibility
StartLimitIntervalSec=10min
StartLimitBurst=3

[Service]
Type=exec
SuccessExitStatus=75
TimeoutStartSec=3600
ExecStart=/usr/bin/flock -E 75 -w 900 %t/topgrade-system.lock /usr/bin/topgrade -cy --no-ask-retry --auto-retry 2 --notify-end on_failure --only system
Nice=19
IOSchedulingClass=idle
CPUSchedulingPolicy=idle
Restart=on-failure
RestartSec=30
