# files/system/etc/systemd/system/btrfs-maintenance.timer
[Unit]
Description=Monthly BTRFS Scrub and Balance

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=timers.target