# files/system/usr/lib/systemd/system/rpm-ostreed-automatic.service.d/override.conf
[Unit]
Wants=network-online.target
After=network-online.target

[Service]
# Run with low scheduling pressure to prefer idle periods.
Nice=19
IOSchedulingClass=idle
