# files/system/usr/lib/systemd/user/flatpak-user-update.service.d/override.conf
[Unit]
Wants=network-online.target
After=network-online.target

[Service]
# Fail fast without network and prefer idle runtime behavior.
ExecStartPre=/usr/bin/nm-online -q --timeout=60
Nice=19
IOSchedulingClass=idle
