# files/system/usr/lib/systemd/system/rpm-ostreed-automatic.timer.d/60-update-schedule.conf
[Unit]
Description=rpm-ostree Automatic Update Trigger

[Timer]
OnCalendar=
OnCalendar=*-*-* 0/4:00:00
OnBootSec=10m
RandomizedDelaySec=5m
