# files/system/usr/lib/systemd/system/rpm-ostreed-automatic.timer.d/override.conf
[Unit]
Description=rpm-ostree Automatic Update Trigger (Hourly)

[Timer]
# Clear default timers
OnBootSec=
OnUnitActiveSec=
OnCalendar=
# Run every hour with a slight random delay to avoid network spikes
OnCalendar=hourly
RandomizedDelaySec=10m
Persistent=true