# files/system/usr/lib/systemd/system/rpm-ostreed-automatic.timer.d/override.conf
[Unit]
Description=rpm-ostree Automatic Update Trigger (45 min)

[Timer]
# Wait for post-boot stabilization, then check every 45 minutes.
OnBootSec=10m
OnUnitActiveSec=45m
RandomizedDelaySec=5m
Persistent=true
AccuracySec=1m
